# scripts/asset-management.service
# Systemd service file for Asset Management System
# Copy this file to /etc/systemd/system/asset-management.service
# Then run: sudo systemctl enable asset-management && sudo systemctl start asset-management

[Unit]
Description=Asset Management System
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/asset-management
ExecStart=/usr/local/bin/docker-compose -f docker-compose.prod.yml up -d
ExecStop=/usr/local/bin/docker-compose -f docker-compose.prod.yml down
ExecReload=/usr/local/bin/docker-compose -f docker-compose.prod.yml restart
TimeoutStartSec=0
User=root

[Install]
WantedBy=multi-user.target